# supervisor_pdks.conf
[program:pdks-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pdks/artisan queue:work database --sleep=3 --tries=3 --max-time=3600 --backoff=10
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
# 10.000 Firma için paralel worker sayısı: 10
numprocs=10
redirect_stderr=true
stdout_logfile=/var/www/pdks/storage/logs/worker.log
stopwaitsecs=3600

[program:pdks-failed-jobs]
# Her başarısız olan job u 3 dk sonra tekrar dene, sonra sil
command=php /var/www/pdks/artisan queue:retry all
autostart=false
# Not: Failed jobs otomasyonu genelde Cron'dan çalışır ama supervisor'la da ayarlanabilir.

[program:pdks-reverb]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/pdks/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/pdks/storage/logs/reverb.log
